# cookbooks/kuansim/templates/default/kuansim.conf.erb
description     "kuansim"

start on (startup
          and filesystem)
task
script
    exec pgrest --db tcp://kuansim:password@localhost/kuansim
end script

